# wechat-autopost.service
[Unit]
Description=WeChat Auto Post Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/.openclaw/workspace/wechat-autopost
ExecStart=/usr/bin/python3 /root/.openclaw/workspace/wechat-autopost/main.py
Restart=always
RestartSec=10

# 日志
StandardOutput=append:/root/.openclaw/workspace/wechat-autopost/logs/service.log
StandardError=append:/root/.openclaw/workspace/wechat-autopost/logs/error.log

[Install]
WantedBy=multi-user.target
